# crontab
# run command once per hour for testing.
1 * * * * /backup.sh >> /var/log/cron.log 2>&1
